# tailord/tailord.service.in
[Unit]
Description=Tux Tailor hardware control service
After=systemd-logind.service

[Service]
Type=dbus
BusName=com.tux.Tailor
ExecStart=@BIN@
Environment="RUST_BACKTRACE=1"

[Install]
WantedBy=multi-user.target